# systemd/webpage-watcher.timer
[Unit]
Description=WebpageWatcher alle 15 Minuten

[Timer]
OnCalendar=*-*-* *:00,15,30,45:00
Persistent=true

[Install]
WantedBy=timers.target
