# systemd/webpage-watcher.service
[Unit]
Description=WebpageWatcher – Seiten- und Erreichbarkeits-Checks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Installationspfad anpassen:
WorkingDirectory=/opt/webpage-watcher
ExecStart=/opt/webpage-watcher/.venv/bin/python /opt/webpage-watcher/monitor.py -c /opt/webpage-watcher/config.yaml
# Optional: Geheimnisse nur per Umgebung
# EnvironmentFile=/etc/webpage-watcher.env

[Install]
WantedBy=multi-user.target
